# deploy/porthole.service
[Unit]
Description=porthole — coral 助理 web GUI (Fastify + Vite, tailscale)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=<your-project-base>/porthole
# 設定(HOST/PORT/PORTHOLE_BASE)由 repo root 的 .env 提供(server/env.ts 載入)。
# 刻意不在此設 Environment=,否則會壓過 .env(loadEnvFile 不覆寫既有 env)。
# 綁定:.env 設 HOST=<tailscale-ip> → 只 tailnet 可達(非 0.0.0.0,不開 LAN)。
# bash -lc 載入 nvm,讓 node/npm 在 PATH
ExecStart=/bin/bash -lc 'npm run start'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
